# ops/runtime/readle-runtime-verify.service
[Unit]
Description=Readle rootful Podman runtime verifier
Wants=network-online.target podman-restart.service readle-podman-forward.service
After=network-online.target podman-restart.service readle-podman-forward.service

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/readle-runtime-verify
TimeoutStartSec=120s

[Install]
WantedBy=multi-user.target
