# ops/mysql-backup/readle-mysql-backup.service
[Unit]
Description=Readle MySQL logical backup to S3
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/local/libexec/readle-mysql-backup
TimeoutStartSec=3600s

[Install]
WantedBy=multi-user.target
